# tags/version-1.6.1:1.1.5/staticDHCPd/samples/ca.uguu.puukusoft.staticDHCPd.plist
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>ca.uguu.puukusoft.staticDHCPd</string>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/bin/env</string>
            <string>python</string>
            <string>-OO</string>
            <string>/usr/local/bin/staticDHCPd</string>
        </array>
        <key>KeepAlive</key>
        <true/>
        <key>RunAtLoad</key>
        <true/>
        <key>ServiceDescription</key>
        <string>Provides staticDHCPd</string>
    </dict>
</plist>
