# tags/version-1.6.3:1.1.5/staticDHCPd/samples/staticDHCPd
### BEGIN INIT INFO
# Provides:          staticDHCPd
# Required-Start:    $all
# Required-Stop:     $local_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start staticDHCPd when all dependencies are up
# Description:       Provides static DHCP "leases" to clients
### END INIT INFO
#LSB-friendly Debian-style initscript.

case "$1" in
start)
    echo -n "Starting staticDHCPd"
    /usr/bin/env python -OO /usr/local/bin/staticDHCPd &
    echo "."
    ;;

stop)
    echo -n "Stopping staticDHCPd"
    kill -15 `cat /var/run/staticDHCPd.pid`
    echo "."
    ;;

force-reload|restart)
    echo -n "Stopping staticDHCPd"
    kill -15 `cat /var/run/staticDHCPd.pid`
    echo "."
    echo -n "Starting staticDHCPd"
    /usr/bin/env python -OO /usr/local/bin/staticDHCPd &
    echo "."
    ;;

*)
    echo "Usage: /etc/init.d/staticDHCPd {start|stop|restart|force-reload}"
    exit 1
    ;;
esac

exit 0
