# samples/org.homelinux.hamsterx.staticDHCPd.plist
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>org.homelinux.hamsterx.staticDHCPd</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/bin/env</string>
			<string>python2.6</string>
			<string>/opt/hamsterx/staticDHCPd/main.py</string>
		</array>
		<key>KeepAlive</key>
		<true/>
		<key>RunAtLoad</key>
		<true/>
		<key>ServiceDescription</key>
		<string>Provides staticDHCPd</string>
	</dict>
</plist>
